# display-i2c-ha-rp/rootfs/etc/services.d/display/run
#!/usr/bin/with-contenv bash
cd /app
exec /app/venv/bin/python3 display.py
